# .config/runit/sv/pipewire/run
#!/bin/sh
# Redirect lock files to RAM (keeps .config clean)
SVDIR="${XDG_RUNTIME_DIR:-/run/user/$(id -u)}/runit/pipewire"
mkdir -p "$SVDIR"
ln -sT "$SVDIR" ./supervise 2>/dev/null

# Redirect stderr to stdout so svlogd can see the logs
exec 2>&1

# Start the process
exec pipewire
